# resume-site.container
# =============================================================
# resume-site.container — Podman Quadlet unit file
#
# Install:
#   cp resume-site.container ~/.config/containers/systemd/
#   systemctl --user daemon-reload
#   systemctl --user start resume-site
#   systemctl --user enable resume-site
#
# Requires config.yaml at ~/resume-site/config.yaml
# =============================================================

[Unit]
Description=resume-site portfolio engine
After=network-online.target

[Container]
# Phase 35: shipped example pins to the immutable v0.3.1 version tag,
# not the moving ``:latest`` alias. Pinning means a `podman auto-update`
# run will NOT roll you forward to a new release without your explicit
# action — which is the safe default. To opt into rolling auto-upgrades,
# switch to ``ghcr.io/kit3713/resume-site:latest`` and uncomment the
# AutoUpdate line below. For maximum reproducibility, pin to the digest
# form (``ghcr.io/kit3713/resume-site@sha256:...``) — see
# docs/PRODUCTION.md §3.0 for the verify-then-pin workflow.
Image=ghcr.io/kit3713/resume-site:v0.3.1
ContainerName=resume-site
# AutoUpdate=registry    # only meaningful with a moving tag (:latest)

# Phase 22.5: bind to 127.0.0.1 only. The container is reachable solely
# via a reverse proxy on the host (Caddy / nginx / Cloudflare Tunnel).
# The previous ``PublishPort=8080:8080`` bound to every interface, which
# — combined with the XFF trust model #16 ships with — would allow a
# remote attacker to forge ``X-Forwarded-For: 127.0.0.1`` and bypass
# the admin IP gate. See docs/PRODUCTION.md for the reverse-proxy
# callout.
PublishPort=127.0.0.1:8080:8080

Volume=%h/resume-site/config.yaml:/app/config.yaml:ro,Z
Volume=resume-site-data:/app/data:Z
Volume=resume-site-photos:/app/photos:Z
Volume=resume-site-backups:/app/backups:Z

# Direct `manage.py backup` archives at the mounted volume so they
# persist across container rebuilds. The Phase 17.2 systemd timer
# (resume-site-backup.timer) reads from the same path on the host.
Environment=RESUME_SITE_BACKUP_DIR=/app/backups

# Liveness probe: /healthz is intentionally lightweight (no I/O, no DB)
# so a transient SQLite write lock or slow disk can't make systemd kill
# the container. The richer ``/readyz`` endpoint (Phase 21.2) checks
# DB connectivity, migration freshness, photo dir writability, and
# disk headroom — appropriate for orchestrator readiness probes
# (Kubernetes / Nomad), but too aggressive for a Podman healthcheck
# whose only response to failure is to mark the container unhealthy.
HealthCmd=curl -f http://localhost:8080/healthz || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=10s

[Service]
Restart=on-failure
TimeoutStartSec=30

[Install]
WantedBy=default.target
