# resume-site-backup.service
# =============================================================
# resume-site-backup.service — Phase 17.2
#
# One-shot systemd unit that wraps `manage.py backup --prune --keep N`
# against a running resume-site container. Triggered by the matching
# resume-site-backup.timer (daily at 02:00 by default).
#
# Install (rootless / per-user):
#   cp resume-site-backup.service ~/.config/systemd/user/
#   cp resume-site-backup.timer   ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now resume-site-backup.timer
#
# Install (system-wide):
#   sudo cp resume-site-backup.{service,timer} /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now resume-site-backup.timer
#
# Override the retention count without editing the file:
#   systemctl --user edit resume-site-backup.service
#   # then in the editor:
#   #   [Service]
#   #   Environment=RESUME_SITE_KEEP=14
#
# Inspect the most recent run:
#   systemctl --user status resume-site-backup.service
#   journalctl --user -u resume-site-backup.service -n 50
# =============================================================

[Unit]
Description=resume-site nightly backup
# Wait for the resume-site container to be running so `podman exec`
# doesn't race on a still-starting service.
After=resume-site.service
Wants=resume-site.service

[Service]
Type=oneshot

# Default retention: keep the 7 newest archives. Override per host via
# `systemctl edit` (see header) without touching this file.
Environment=RESUME_SITE_KEEP=7

# `podman exec` runs against the long-running resume-site container
# (Quadlet creates it as `resume-site`; compose.yaml uses the same
# container_name). The CLI inside the container resolves the output
# directory via RESUME_SITE_BACKUP_DIR (set on the container itself in
# the Quadlet / compose env block, default /app/backups), so the
# archive lands on the resume-site-backups volume.
ExecStart=/usr/bin/podman exec resume-site \
    python manage.py backup --prune --keep ${RESUME_SITE_KEEP}

# Cap a runaway backup at 30 minutes so the timer keeps to schedule.
TimeoutStartSec=30min

# Don't restart on failure — the timer's Persistent=true catches the
# next missed window. Restart loops would mask real underlying issues
# (disk full, container down).
Restart=no

# -------------------------------------------------------------
# Phase 28.4 — systemd.exec hardening (#27)
#
# This unit shells out to `podman exec` against the long-running
# container. The actual write to the backups volume happens *inside*
# the container; the host-side process is just the podman CLI plus
# its socket-talker. Hardening here closes the host-side surface —
# it does NOT need to whitelist the backup directory, because that
# directory is the container's view, not a host path systemd sees.
#
# Rollback for any single directive: comment out the line and run
# `systemctl --user daemon-reload` then restart the timer.
# -------------------------------------------------------------

# NoNewPrivileges: blocks setuid/setcap escalation. Pure restriction;
# rarely breaks anything. Rollback: comment out the line.
NoNewPrivileges=yes

# PrivateTmp: gives the unit a private /tmp and /var/tmp namespace.
# Pure restriction. Rollback: comment out.
PrivateTmp=yes

# ProtectSystem=strict: mounts /, /usr, /boot, /efi, /etc read-only.
# `podman exec` only needs to read from /usr (the binary) and talk
# to the Podman socket under XDG_RUNTIME_DIR (covered by the implicit
# allowlist that ProtectSystem leaves open). Rollback: drop to
# ProtectSystem=full if a future change writes to /etc.
ProtectSystem=strict

# ProtectHome: hides /home, /root, /run/user. `podman exec` reaches
# the rootless daemon via the socket under %t (XDG_RUNTIME_DIR),
# whitelisted below — it doesn't need /home for this command.
# Rollback: ProtectHome=read-only if a future invocation needs to
# read a config from %h.
ProtectHome=yes

# ReadWritePaths: the Podman socket lives under %t (XDG_RUNTIME_DIR);
# whitelist it so `podman exec` can talk to the daemon under
# ProtectSystem=strict.
ReadWritePaths=%t/containers

# RestrictSUIDSGID: forbids creating setuid/setgid files. Pure
# restriction. Rollback: comment out.
RestrictSUIDSGID=yes

# LockPersonality: locks personality(2) to the boot default. Pure
# restriction. Rollback: comment out.
LockPersonality=yes

# RestrictNamespaces: forbids the unit from creating new namespaces.
# `podman exec` enters an existing container's namespaces; it does
# not create new ones from this process. Rollback: if podman errors
# with EPERM on unshare/clone after a Podman upgrade, set to
# RestrictNamespaces=user mnt or comment out.
RestrictNamespaces=yes

# SystemCallArchitectures=native: blocks non-native syscalls. Pure
# restriction on a 64-bit host. Rollback: comment out.
SystemCallArchitectures=native

# MemoryDenyWriteExecute=yes  # NOTE: untested against Pillow's image-processing
                              # path. Enable in a staging deployment first; if
                              # Pillow segfaults on photo upload, comment out.

[Install]
# Pulled in as a dependency by the timer; no separate WantedBy needed.
WantedBy=default.target
